# etc/init/bloonix-wtrm.in
#! /bin/sh

# Red Hat style
# chkconfig: 2345 99 10
# description: bloonix wtrm

# Debian style
### BEGIN INIT INFO
# Provides:          bloonix-wtrm
# Required-Start:    $network $syslog $all
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the bloonix-server daemon.
# Description:       bloonix-server is a program that communicates
#                    with the bloonix server.
### END INIT INFO

. @@USRLIBDIR@@/bloonix/bin/bloonix-init-source
CHECK_DIRS="$CACHEDIR $LOGDIR $RUNDIR"
__run
